# server/systemd/local-chat-proxy.service
[Unit]
Description=LC Chat Transparent HTTP Proxy
After=network.target local-chat-server.service
Wants=local-chat-server.service

[Service]
Type=simple
User=root
ExecStart=/home/gamecat/lc-server/venv/bin/python /home/gamecat/lc-server/app/services/http_proxy.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
